# docs/migration/launchagents/com.linze.openai-database.codex-daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.linze.openai-database.codex-daily-sync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/linzezhang/Library/Application Support/OpenAIDatabase/MemoryAtlas/github-main/OpenAIDatabase/scripts/run_codex_memory_auto_update.py</string>
		<string>--database-dir</string>
		<string>/Users/linzezhang/Library/Application Support/OpenAIDatabase/MemoryAtlas/github-main/OpenAIDatabase</string>
		<string>--publish-runtime</string>
		<string>--commit</string>
		<string>--push</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/linzezhang/Library/Logs/OpenAIDatabase/codex-daily-sync.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/linzezhang/Library/Logs/OpenAIDatabase/codex-daily-sync.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Hour</key>
			<integer>3</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>1</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>3</integer>
			<key>Minute</key>
			<integer>0</integer>
			<key>Weekday</key>
			<integer>5</integer>
		</dict>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/linzezhang/Library/Application Support/OpenAIDatabase/MemoryAtlas/github-main/OpenAIDatabase</string>
</dict>
</plist>
